# start.ini
[supervisord]
nodaemon=true


[program:apache2]
command=/usr/sbin/apache2ctl -D "FOREGROUND" -k start
directory=/
user=root
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true
stopasgroup=true
stdout_logfile=/app/logs/apache2out.log
stderr_logfile=/app/logs/apache2err.log
